# reports/crontab
# Каждый день в 9:00 UTC
0 9 * * * root . /etc/environment && /usr/local/bin/python /app/daily_report.py >> /var/log/cron.log 2>&1
